# Manager/supervisor-app.ini
[program:app-uwsgi]
command = /sbin/uwsgi --http-socket :80 --plugin python --wsgi-file /var/lib/gncloud/Manager/__init__.py --logto /var/log/gncloud/manager.log --processes %(ENV_WSGI_PROCESSES)s --threads %(ENV_WSGI_THREADS)s --callable app